# dns/supervisord.conf
[supervisord]
nodaemon=true

[program:restartdns]
command=/restartdns.sh
numprocs=1
autostart=true
autorestart=true

[program:dnsmasq]
command=dnsmasq -d
numprocs=1
autostart=true
autorestart=true
