# node/supervisord.conf
[supervisord]
nodaemon=true

[program:neo4j]
command=/var/lib/neo4j/bin/neo4j console
startsecs=40
numprocs=1
autostart=true
autorestart=true
